# deploy/giftbot.service
[Unit]
Description=Gift Bot — Telegram gifts (multi-user, optional)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/giftbot
ExecStart=/usr/bin/python3 /root/giftbot/main.py
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
